# cho_v5/in4_landscape/services--c/virtualization--o/vm--f/xen/rev5/5_service/profile.d/rev5__services--c--virtualization--o--vm--f--xen--g--main--s.aliases.sh
alias rev5__services--c--virtualization--o--vm--f--xen--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__xen--g__sdo-...=FULL_ALIAS

alias rev5__services--c--virtualization--o--vm--f--xen--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__xen--g__do-...=FULL_ALIAS


