# cho_v5/in4_landscape/internals--c/management--o/rev5--f/core/rev5/5_service/profile.d/rev5__internals--c--management--o--rev5--f--core--g--main--s.aliases.sh
alias rev5__internals--c--management--o--rev5--f--core--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__core--g__sdo-...=FULL_ALIAS

alias rev5__internals--c--management--o--rev5--f--core--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__core--g__do-...=FULL_ALIAS


