# cho_v5/_in4_landscape_old/internals--c/linux_sys--o/boot--f/grub2/rev5/5_service/profile.d/rev5__internals--c--linux_sys--o--boot--f--grub2--g--main--s.aliases.sh
alias rev5__internals--c--linux_sys--o--boot--f--grub2--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__grub2--g__sdo-...=FULL_ALIAS

alias rev5__internals--c--linux_sys--o--boot--f--grub2--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__grub2--g__do-...=FULL_ALIAS


