# cho_v4/services--c/vpn:o/ssl:f/ocserv:g/_systemd/in4__openconnect@.socket
[Unit]
Description=in4 OpenConnect VPN socket - %i
Wants=network.target

[Socket]
ListenFIFO=/run/openconnect/%i
PipeSize=128
