# cho_v3/ontology/security/sssd/engine/profile.d/power.sh
alias srm='sudo /usr/bin/rm'
alias smcedit='sudo /usr/bin/mcedit'
alias svi='sudo /usr/bin/vi'