# cho_v5/in4_landscape/data_safety--c/replication--o/block--f/drbd9/rev5/5_service/profile.d/rev5__data_safety--c--replication--o--block--f--drbd9--g--main--s.aliases.sh
alias rev5__data_safety--c--replication--o--block--f--drbd9--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__drbd9--g__sdo-...=FULL_ALIAS

alias rev5__data_safety--c--replication--o--block--f--drbd9--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__drbd9--g__do-...=FULL_ALIAS


