# cho_v5/in4_landscape/data_safety--c/replication--o/block--f/drbd9/rev5/5_service/systemd/rev5__data_safety--c--replication--o--block--f--drbd9--g--primary--s.socket
[Socket]
ListenStream=80
KeepAlive=yes
NoDelay=yes
IPTTL=10
IPTOS=low-delay
TCPCongestion=cubic

#ExecStartPre=
#ExecStartPost=
#ExecStopPre=
#ExecStopPost=

[Install]
WantedBy=sockets.target
