# cho_v4/data_safety:c/snapshot:o/others:f/btrfs--g/service/in4__snapshot-app_hourly.timer
[Unit]
Description=Timer for in4__snapshot
Requires=var-run.mount nss-lookup.target network.target remote-fs.target
After=var-run.mount nss-lookup.target network.target remote-fs.target

[Install]
WantedBy=multi-user.target

[Timer]
OnUnitActiveSec=1h
AccuracySec=10m
Unit=in4__snapshot-app.service 
