# cho_v5/_in4_landscape_old/data_safety--c/file_backup--o/infrastructure--f/bareos/rev5/5_service/profile.d/rev5__data_safety--c--file_backup--o--infrastructure--f--bareos--g--main--s.aliases.sh
alias rev5__data_safety--c--file_backup--o--infrastructure--f--bareos--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__bareos--g__sdo-...=FULL_ALIAS

alias rev5__data_safety--c--file_backup--o--infrastructure--f--bareos--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__bareos--g__do-...=FULL_ALIAS


