# cho_v5/in4_landscape/internals--c/management--o/in4/in4/5_service/profile.d/in4__main--s.aliases.sh
########    #######    ########    #######    ########    ########
##     / / / /    License    \ \ \ \ 
##    Copyleft culture, Copyright (C) is prohibited here
##    This work is licensed under a CC BY-SA 4.0
##    Creative Commons Attribution-ShareAlike 4.0 License
##    Refer to the http://creativecommons.org/licenses/by-sa/4.0/
########    #######    ########    #######    ########    ########
##    / / / /    Code Climate    \ \ \ \ 
##    Language = bash DSL, profiles
##    Indent = space;    4 chars;
########    #######    ########    #######    ########    ########

alias in4__class--c--order--o--family--f--genus--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias in4__short__genus--g__sdo-...=FULL_ALIAS

alias in4__class--c--order--o--family--f--genus--g--main--s__do-...='FULL_PATH ARGS ...'
alias in4__short__genus--g__do-...=FULL_ALIAS


