# cho_v5/_in4_landscape_old/services--c/webui--o/analytics--f/kibana4/rev5/1_build/obs/kibana4/kibana.service
#
# kibana systemd service
#

[Unit]
Description=Kibana
Documentation=https://www.elastic.co/guide
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/kibana
User=kibana
Group=kibana
ExecStart=/opt/kibana/bin/kibana

[Install]
WantedBy=multi-user.target

