# cho_v5/_in4_landscape_old/internals--c/linux_sys--o/kernel_leap42_2--f/kernel-source/rev5/5_service/profile.d/rev5__internals--c--linux_sys--o--kernel_leap42_2--f--kernel-source--g--main--s.aliases.sh
########    #######    ########    #######    ########    ########
##     / / / /    License    \ \ \ \ 
##    Copyleft culture, Copyright (C) is prohibited here
##    This work is licensed under a CC BY-SA 4.0
##    Creative Commons Attribution-ShareAlike 4.0 License
##    Refer to the http://creativecommons.org/licenses/by-sa/4.0/
########    #######    ########    #######    ########    ########
##    / / / /    Code Climate    \ \ \ \ 
##    Language = bash DSL, profiles
##    Indent = space;    4 chars;
########    #######    ########    #######    ########    ########

alias rev5__internals--c--linux_sys--o--kernel_leap42_2--f--kernel-source--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__kernel-source--g__sdo-...=FULL_ALIAS

alias rev5__internals--c--linux_sys--o--kernel_leap42_2--f--kernel-source--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__kernel-source--g__do-...=FULL_ALIAS


