# cho_v5/in4_landscape/logitoring--c/messagebus--o/syslog--f/rsyslog/rev5/5_service/profile.d/rev5__logitoring--c--messagebus--o--syslog--f--rsyslog--g--main--s.aliases.sh
alias rev5__logitoring--c--messagebus--o--syslog--f--rsyslog--g--main--s__sdo-...='sudo FULL_PATH ARGS ...'
alias rev5__short__rsyslog--g__sdo-...=FULL_ALIAS

alias rev5__logitoring--c--messagebus--o--syslog--f--rsyslog--g--main--s__do-...='FULL_PATH ARGS ...'
alias rev5__short__rsyslog--g__do-...=FULL_ALIAS


